# personal/luci-app-xraycustom/root/etc/uci-defaults/99-xray
#!/bin/sh

chmod +x /usr/bin/xraycustom
ls -l /usr/bin/xraycustom